# gutentag.cron
# Gutentag, World! — runs every minute and appends to log
* * * * * echo "Gutentag, World!" >> /tmp/gutentag.log
